# provisioning/roles/utils/cron/templates/etc/systemd/system/cron.timer.j2
{{ ansible_managed|comment }}

[Unit]
Description={{ cron_timer_description }}

[Timer]
OnCalendar={{ cron_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
