# provisioning/roles/web/templates/compute-ranking.timer.j2
[Unit]
Description=Compute ranking daily
After=postgresql
Type=oneshot
ExecStart={{ manage_path }} ranking

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
